# setup_files/docker/supervisor/supervisord.conf
; docker/supervisor/supervisord.conf
;
; Supervisor runs inside the app container and starts Laravel Octane.
; Queue workers and scheduler run as SEPARATE containers (see docker-compose.yml),
; so this file only manages the Octane HTTP server.

[supervisord]
nodaemon=true
user=www
logfile=/var/www/html/storage/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── Laravel Octane (Swoole HTTP server) ──────────────────────────────────────
[program:octane]
command=php /var/www/html/artisan octane:start
    --server=swoole
    --host=0.0.0.0
    --port=8000
    --workers=%(ENV_OCTANE_WORKERS)s
    --task-workers=%(ENV_OCTANE_TASK_WORKERS)s
    --max-requests=%(ENV_OCTANE_MAX_REQUESTS)s

directory=/var/www/html
user=www
numprocs=1
autostart=true
autorestart=true
startretries=5
startsecs=5

; Stream logs to Docker's stdout/stderr so `docker compose logs -f app` works
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Stop signal — SIGTERM triggers Octane's graceful shutdown
stopsignal=SIGTERM
stopwaitsecs=30

environment=HOME="/home/www",USER="www"
